# deploy/systemd/bond-notify-worker.service
[Unit]
Description=SimpleBond Notification Worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Replace with your deploy user and repo path.
User=simplebond
WorkingDirectory=/opt/simple-bond
ExecStart=/usr/bin/node backend/worker.mjs
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production
Environment=BOND_NOTIFY_BASE_URL=https://api.bond.futarchy.ai
Environment=SIMPLE_BOND_FRONTEND_URL=https://bond.futarchy.ai
EnvironmentFile=-/etc/simple-bond/bond-notify.env

[Install]
WantedBy=multi-user.target
